# deploy/atlantis.service
[Unit]
Description=atlantis gRPC server
Wants=network-online.target
After=network-online.target docker.service postgresql.service memcached.service
Requires=docker.service

[Service]
Restart=on-failure
RestartSec=5
# Reap any container left behind by an unclean stop before starting.
ExecStartPre=-/usr/bin/docker rm -f atlantis
# Image tag hard-coded to atlantis:local. `make deploy` tags the build as
# $(IMAGE) and restarts this unit; if you override IMAGE at build time,
# edit this line and `systemctl daemon-reload` or the restart will pull
# the old tag.
ExecStart=/usr/bin/docker run --rm --name atlantis \
  --network host \
  --env-file /etc/atlantis/atlantis.env \
  -v /etc/atlantis/tls:/etc/atlantis/tls:ro \
  atlantis:local
# atlantis drains on SIGTERM (~30s worst case: 15s gRPC GracefulStop + 10s
# outbox drain + 5s health shutdown). docker run forwards SIGTERM to PID 1;
# give the drain room before systemd sends SIGKILL.
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
